# system_integration/systemd/orin_stack.service.d/override.conf
[Service]
# ---- Environment ----
Environment=COLCON_TRACE=
Environment=AMENT_TRACE_SETUP_FILES=
Environment=AMENT_PYTHON_EXECUTABLE=/usr/bin/python3
Environment=COLCON_PYTHON_EXECUTABLE=/usr/bin/python3
Environment=ROS2CLI_NO_DAEMON=1
Environment=ROS_DOMAIN_ID=0

# Force Fast DDS
Environment=RMW_IMPLEMENTATION=rmw_fastrtps_cpp
Environment=FASTRTPS_DEFAULT_PROFILES_FILE=/home/sauce/.ros/fastdds.xml
Environment=CYCLONEDDS_URI=

# ---- Logging ----
SyslogIdentifier=orin_stack
StandardOutput=journal
StandardError=journal

# ---- Clean start ----
ExecStartPre=
ExecStartPre=/bin/bash -lc 'echo "[orin_stack][pre] $(date -Is) domain=$ROS_DOMAIN_ID rmw=$RMW_IMPLEMENTATION"'

# ---- Main launch ----
ExecStart=
ExecStart=/bin/bash -lc 'set +u; source /opt/ros/humble/setup.bash; \
  [ -f /data/ros_ws/install/setup.bash ] && source /data/ros_ws/install/setup.bash; \
  exec /data/orin_stack'

# ---- Clean stop (never fail the unit) ----
ExecStopPost=
ExecStopPost=/bin/bash -lc '/usr/local/bin/orin.oak stop; rc=$?; echo "[orin_stack][stoppost] orin.oak rc=$rc"; exit 0'

# ---- Shutdown behavior ----
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=5
